# conf/service/system-monitor.service
[Unit]
Description=System Monitoring Service

[Service]
Type=simple
ExecStart=/usr/local/sbin/system-monitor
Restart=on-failure
