# conf/crontab/caemlyn.cron
*/30  * * * * tcsh -c "$HOME/bin/proc-mem -k 50 thunderbird"
