# conf/crontab/cibola.cron
13    0 * * * updatedb -l 0 -o $HOME/local/var/mlocate.db -U $HOME
