# conf/crontab/zadash.cron
*     * * * * xdotool search --name URGENT >/dev/null || sudo -H -u buddy -i $HOME/bin/monitor-pidgin-unseen caemlyn >/dev/null 2>&1
